# deploy/systemd/tour-guide-frontend.service
[Unit]
Description=Tour Guide Frontend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/frontend
EnvironmentFile=-/etc/default/tour-guide-frontend
ExecStartPre=/usr/bin/test -f /opt/frontend/server.js
ExecStart=/usr/bin/node /opt/frontend/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=API_URL=https://yourdomain.com
StandardOutput=journal
StandardError=journal
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
